# im2modelpy_api/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=python3 -u /app/app.py
autostart=true
autorestart=true
